# deploy/supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid
childlogdir=/var/log/supervisor

# ---- MySQL ----
[program:mysql]
command=/usr/sbin/mysqld --user=mysql --datadir=/var/lib/mysql
autostart=true
autorestart=true
priority=10
stdout_logfile=/var/log/supervisor/mysql.log
stderr_logfile=/var/log/supervisor/mysql-err.log

# ---- MinIO ----
[program:minio]
command=/usr/local/bin/minio server /data/minio --console-address ":9001"
environment=MINIO_ROOT_USER="minioadmin",MINIO_ROOT_PASSWORD="minioadmin"
autostart=true
autorestart=true
priority=20
stdout_logfile=/var/log/supervisor/minio.log
stderr_logfile=/var/log/supervisor/minio-err.log

# ---- Backend (FastAPI + Milvus Lite 内嵌) ----
[program:backend]
command=python3.11 -m uvicorn app.main:app --host 0.0.0.0 --port 8080
directory=/app/backend
environment=SECRET_KEY="change-me-in-production",MYSQL_HOST="127.0.0.1",MYSQL_PORT="3306",MYSQL_USER="root",MYSQL_PASSWORD="root",MYSQL_DATABASE="museum_monitor",MILVUS_URI="/app/backend/data/milvus_lite.db",MINIO_ENDPOINT="127.0.0.1:9000",MINIO_ACCESS_KEY="minioadmin",MINIO_SECRET_KEY="minioadmin"
autostart=true
autorestart=true
priority=30
startsecs=5
stdout_logfile=/var/log/supervisor/backend.log
stderr_logfile=/var/log/supervisor/backend-err.log

# ---- Nginx (Frontend) ----
[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
autostart=true
autorestart=true
priority=40
stdout_logfile=/var/log/supervisor/nginx.log
stderr_logfile=/var/log/supervisor/nginx-err.log
